# fgn-backup.crontab
# install this file (as fgn-backup, if you like) into /etc/cron.d
# m h dom mon dow user command
  2 2  *   *   *  root test -x /usr/local/bin/fgn-backup && /usr/local/bin/fgn-backup
